# systemd/agentic-crm.service
[Unit]
Description=Agentic CRM - AI agents for media ad sales
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/claude/crm-azteca
ExecStart=/usr/bin/node engine/dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=/root/claude/crm-azteca/.env

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/root/claude/crm-azteca/store /root/claude/crm-azteca/data /root/claude/crm-azteca/groups

[Install]
WantedBy=default.target
